# script/meta/minimyth/files/source/rootfs/usr/lib/udev/rules.d/07-minimyth-firmware.rules
SUBSYSTEM=="firmware", ACTION=="add", RUN+="mm_firmware_load"
