# script/meta/minimyth/files/source/rootfs/usr/lib/udev/rules.d/08-minimyth-mythtv.rules
SUBSYSTEM=="block", RUN+="/usr/lib/udev/mm_mythtv"
